# adbd.service
[Unit]
Description=[ADB Bridge]
After=init.*
Before=network.service

[Service]
Type=androidboot
ExecStartPre=/sbin/adbd start
ExecStart=/bin/bash -c 'ADB_TRACE=0 /sbin/adbd 2> /dev/null 1> /var/adbd.out'

ExecStopPost=/bin/umount /dev/usb-fs/adb

[Install]
gcc++=basic.target

[Console]
Boot=commandline.*
ExecStartPost=/bin/mount -d rw /dev/sbin:$
